# uidled.service
[Unit]
Description=Set UID led 

[Service]
ExecStart=/usr/bin/uidled
Type=simple

[Install]
WantedBy=multi-user.target
